# 136_/var/lib/opkg/info/udev.list
/lib/udev/keymaps/hewlett-packard
/lib/udev/keymaps/module-asus-w3j
/lib/udev/keymaps/micro-star
/lib/udev/keymaps/module-sony
/lib/udev/keymaps/force-release/common-volume-keys
/lib/udev/keymaps/module-lenovo
/lib/udev/scsi_id
/lib/udev/keymaps/lenovo-thinkpad_x6_tablet
/lib/udev/keymaps/samsung-sq1us
/lib/udev/keymaps/fujitsu-esprimo_mobile_v5
/lib/udev/keymaps/fujitsu-esprimo_mobile_v6
/lib/udev/keymaps/oqo-model2
/lib/udev/devices/
/lib/udev/cdrom_id
/lib/udev/keymap
/lib/udev/keymaps/
/lib/udev/keymaps/module-sony-old
/lib/udev/keymaps/acer-aspire_5720
/lib/udev/mobile-action-modeswitch
/lib/udev/keymaps/hewlett-packard-presario-2100
/lib/udev/keymaps/maxdata-pro_7000
/lib/udev/pci-db
/etc/udev/rules.d/local.rules
/lib/udev/keymaps/hewlett-packard-pavilion
/lib/udev/keymaps/fujitsu-amilo_pa_2548
/lib/udev/keymaps/zepto-znote
/lib/udev/path_id
/lib/udev/hid2hci
/lib/udev/keymaps/medion-fid2060
/lib/udev/keymaps/lenovo-ideapad
/lib/udev/keymaps/hewlett-packard-tx2
/lib/udev/keymaps/asus
/lib/udev/keymaps/onkyo
/lib/udev/usb-db
/lib/udev/keymaps/force-release/samsung-other
/lib/udev/keymaps/dell
/lib/udev/keymaps/lenovo-thinkpad-usb-keyboard-trackpoint
/lib/udev/keymaps/acer-aspire_6920
/etc/udev/
/lib/udev/keymaps/lenovo-3000
/lib/udev/keymaps/everex-xt5000
/lib/udev/keymaps/fujitsu-amilo_pro_v3205
/sbin/udevd
/lib/udev/keymaps/acer-aspire_5920g
/lib/udev/edd_id
/lib/udev/input_id
/lib/udev/keymaps/inventec-symphony_6.0_7.0
/lib/udev/keymaps/hewlett-packard-2510p_2530p
/lib/udev/keymaps/toshiba-satellite_a100
/etc/udev/scripts/network.sh
/lib/udev/keymaps/toshiba-satellite_a110
/lib/udev/keymaps/medionnb-a555
/lib/udev/keymaps/acer-travelmate_c300
/lib/udev/mtd_probe
/lib/udev/
/lib/udev/ata_id
/lib/udev/keymaps/logitech-wave-cordless
/lib/udev/collect
/lib/udev/keymaps/module-ibm
/lib/udev/findkeyboards
/lib/udev/keymaps/samsung-sx20s
/lib/udev/keyboard-force-release.sh
/lib/udev/keymaps/fujitsu-amilo_pro_edition_v3505
/etc/init.d/udev
/lib/udev/create_floppy_devices
/lib/udev/usb_id
/lib/udev/write_cd_rules
/lib/udev/keymaps/olpc-xo
/lib/udev/rules.d/
/lib/udev/rules.d/50-firmware.rules
/lib/udev/keymaps/logitech-wave
/etc/udev/scripts/
/lib/udev/keymaps/hewlett-packard-tablet
/lib/udev/keymaps/acer
/etc/udev/udev.conf
/lib/udev/keymaps/force-release/
/lib/udev/keymaps/logitech-wave-pro-cordless
/etc/udev/rules.d/
/lib/udev/keymaps/lg-x110
/etc/udev/scripts/mount.sh
/lib/udev/keymaps/toshiba-satellite_m30x
/lib/udev/keymaps/ibm-thinkpad-usb-keyboard-trackpoint
/lib/udev/v4l_id
/lib/udev/fstab_import
/lib/udev/keymaps/samsung-other
/lib/udev/firmware
/lib/udev/keymaps/hewlett-packard-compaq_elitebook
/lib/udev/keymaps/fujitsu-amilo_si_1520
/lib/udev/rule_generator.functions
/lib/udev/keymaps/lenovo-thinkpad_x200_tablet
/lib/udev/keymaps/compaq-e_evo
/lib/udev/write_net_rules
