# 135_/var/lib/opkg/info/atfwd-daemon.postrm
#!/bin/sh
update-rc.d $D atfwd-daemon remove
