# 135_/var/lib/opkg/info/csd-server.list
/etc/init.d/csdserver
/usr/bin/csd_server
