# 135_/var/lib/opkg/info/csd-server.postrm
#!/bin/sh
update-rc.d $D csdserver remove
