# 135_/var/lib/opkg/info/data.list
/tmp/tests/dss_test_netc
/usr/lib/libqcmap_cm.so.1
/usr/bin/netmgrd
/usr/lib/libdssnew.so.1
/tmp/tests/
/usr/bin/radish
/etc/init.d/qti
/usr/lib/libdsnet.so.1
/etc/mobileap_cfg.xml
/etc/mobileap_cfg.xsd
/usr/lib/libdsstest.so.1
/etc/init.d/netmgrd
/etc/init.d/data-init
/etc/init.d/QCMAP_ConnectionManager
/etc/udhcpc.d/udhcpc.script
/usr/bin/QCMAP_ConnectionManager
/usr/bin/qti
/tmp/tests/dssnetc.sh
/usr/lib/libnetmgr.so.0
