# 135_/var/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.d/
/etc/dnsmasq.conf
/usr/bin/dnsmasq
/etc/init.d/dnsmasq
