# 135_/var/lib/opkg/info/modem-shutdown.postinst
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT="-s"
fi
update-rc.d $OPT modem-shutdown stop 01 0 6 .
