# 135_/var/lib/opkg/info/modem-shutdown.postrm
#!/bin/sh
update-rc.d $D modem-shutdown remove
