# 135_/var/lib/opkg/info/modem-shutdown.list
/bin/modem-shutdown
/bin/
/etc/init.d/modem-shutdown
