# 135_/var/lib/opkg/info/qmi-shutdown-modem.postrm
#!/bin/sh
update-rc.d $D qmi_shutdown_modemd remove
