# 135_/var/lib/opkg/info/qmi-shutdown-modem.list
/usr/
/etc/
/etc/init.d/qmi_shutdown_modemd
/usr/bin/
/usr/bin/qmi_shutdown_modem
/etc/init.d/
/
