# 135_/var/lib/opkg/info/sysvinit.list
/sbin/poweroff.sysvinit
/usr/bin/last.sysvinit
/etc/rc2.d/S99stop-bootlogd
/etc/default/rcS
/usr/bin/mesg.sysvinit
/etc/rcS-zte
/etc/rc4.d/S99stop-bootlogd
/sbin/reboot.sysvinit
/etc/init.d/rc
/bin/mountpoint
/sbin/bootlogd
/etc/init.d/bootlogd
/etc/init.d/stop-bootlogd
/sbin/telinit
/etc/rc3.d/S99stop-bootlogd
/etc/rcS.d/S07bootlogd
/sbin/fstab-decode
/sbin/init.sysvinit
/usr/bin/wall.sysvinit
/etc/init.d/rcS
/sbin/runlevel
/sbin/shutdown.sysvinit
/usr/bin/utmpdump
/sbin/halt.sysvinit
/etc/rcS-other
/etc/rc5.d/S99stop-bootlogd
/usr/bin/lastb
