# 135_/var/lib/opkg/info/time-services.list
/etc/rc0.d/
/etc/rc1.d/
/etc/rc2.d/
/etc/rc3.d/S01time_serviced
/etc/rc3.d/
/etc/rc6.d/K01time_serviced
/etc/rc1.d/K01time_serviced
/etc/rc4.d/
/usr/bin/time_daemon
/etc/init.d/time_serviced
/etc/rc5.d/S01time_serviced
/etc/rc5.d/
/etc/rc6.d/
/etc/rc2.d/S01time_serviced
/etc/rc0.d/K01time_serviced
/etc/rc4.d/S01time_serviced
