# 135_/var/lib/opkg/info/udev-cache.postrm
#!/bin/sh
update-rc.d $D udev-cache remove
